# psi-notify.openrc
#!/sbin/openrc-run
command=/usr/bin/psi-notify

depent(){
need xdm
}
